# launchd/codex-hot-swap.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>__LABEL__</string>

  <key>ProgramArguments</key>
  <array>
    <string>__PREFIX__/codex-predictive-daemon</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>CODEX_HOME</key>
    <string>__CODEX_HOME__</string>
    <key>PATH</key>
    <string>__PATH__</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__CODEX_HOME__/codex-hot-swap-daemon.out.log</string>

  <key>StandardErrorPath</key>
  <string>__CODEX_HOME__/codex-hot-swap-daemon.err.log</string>
</dict>
</plist>
